# archiso/airootfs/etc/systemd/system/aios-model-init.service
[Unit]
Description=AIos - Import Qwen3 0.6B into Ollama
After=network.target
ConditionPathExists=/usr/local/lib/archspeech/models/qwen3-0.6b.gguf
ConditionPathExists=!/var/lib/aios/.model_ready

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/archspeech/bin/python3 /usr/local/bin/aios-model-setup
TimeoutStartSec=600
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
