# archiso/airootfs/etc/systemd/system/ollama.service.d/debug.conf
[Service]
Environment=OLLAMA_MODELS=/tmp/ollama-models
StandardOutput=journal
StandardError=journal
