# archiso/airootfs/etc/systemd/system/ollama.service.d/debug.conf
[Service]
Environment=OLLAMA_DEBUG=1
Environment=OLLAMA_VERBOSE=1
# Use /tmp for writes (live USB cowspace is too small for model blobs)
# Models are pre-staged at /usr/share/ollama/.ollama and copied to /tmp on start
Environment=OLLAMA_MODELS=/tmp/ollama-models
StandardOutput=journal
StandardError=journal
# Copy pre-staged models to writable /tmp location on first start
ExecStartPre=/bin/bash -c 'if [ -d /usr/local/lib/archspeech/ollama/models ] && [ ! -d /tmp/ollama-models ]; then cp -r /usr/local/lib/archspeech/ollama/models /tmp/ollama-models; fi'
