# archiso/airootfs/etc/systemd/system/ollama.service.d/debug.conf
[Service]
Environment=OLLAMA_DEBUG=1
Environment=OLLAMA_VERBOSE=1
StandardOutput=journal
StandardError=journal
