# archiso/airootfs/etc/systemd/system/ollama.service.d/debug.conf
[Service]
Environment=OLLAMA_MODELS=/tmp/ollama-models
StandardOutput=journal
StandardError=journal
# Copy pre-staged models to RAM in background — don't block Ollama startup
ExecStartPre=/bin/bash -c 'mkdir -p /tmp/ollama-models && (cp -r /usr/local/lib/archspeech/ollama/models/. /tmp/ollama-models/ &)'
