# archiso/airootfs/etc/systemd/system/ollama.service.d/debug.conf
[Service]
Environment=OLLAMA_DEBUG=1
Environment=OLLAMA_VERBOSE=1
# Store models in /tmp — RAM-backed, no 256MB cowspace limit on live USB
Environment=OLLAMA_MODELS=/tmp/ollama-models
StandardOutput=journal
StandardError=journal
